# dev/docker/etc/supervisord-mysqld.conf
[program:mysqld]
command=/startscripts/start-mysqld.sh
numprocs=1
autostart=true
autorestart=true